# swifi-theme-bootstrap/root/etc/uci-defaults/30_swifi-theme-bootstrap
#!/bin/sh
uci batch <<-EOF
	set luci.themes.SwifiBootstrap=/luci-static/swifi-bootstrap
	set luci.main.mediaurlbase=/luci-static/swifi-bootstrap
	commit luci
EOF
exit 0
