# swifi-service-coovachilli-3g/files/utils/crontabs/root
* * * * * sh /etc/chilli/utils/chilli_query_list.sh
